# container-files/etc/supervisor.d/gitbucket.conf
[program:gitbucket]
command = /usr/bin/java -jar /opt/gitbucket.war --host=%(ENV_HOSTNAME)s --prefix=%(ENV_CONTEXTPATH)s --gitbucket.home=%(ENV_DATA_DIR)s
autorestart = true
